# docker-compose/supervisor/php-fpm.conf
[program:php-fpm]
user=root
command=/usr/local/sbin/php-fpm
redirect_stderr=true
stdout_logfile=/var/log/php-fpm-supervisor.log